# log-monitor.service
[Unit]
Description=Log Monitor

[Service]
ExecStart=/usr/local/bin/log-monitor.sh
Restart=always

[Install]
WantedBy=multi-user.target
